# deploy/supervisor.conf
[program:buykori-web]
command=/var/www/buykori-adsync/venv/bin/uvicorn app.main:app --host 127.0.0.1 --port 8000 --workers 2 --no-access-log --timeout-keep-alive 75 --limit-max-requests 10000
directory=/var/www/buykori-adsync
user=buykori
autostart=true
autorestart=true
stdout_logfile=/var/log/supervisor/buykori-web.out.log
stderr_logfile=/var/log/supervisor/buykori-web.err.log
stdout_logfile_maxbytes=10MB
stdout_logfile_backups=5
stderr_logfile_maxbytes=10MB
stderr_logfile_backups=5
environment=PATH="/var/www/buykori-adsync/venv/bin:%(ENV_PATH)s",EVENT_INGEST_MODE="redis_stream",USAGE_RESERVATION_MODE="worker",DB_POOL_SIZE="3",DB_MAX_OVERFLOW="3",HTTP_TIMEOUT_SECONDS="8"

[program:buykori-worker]
command=/var/www/buykori-adsync/venv/bin/python -m app.services.event_worker
directory=/var/www/buykori-adsync
user=buykori
autostart=true
autorestart=true
stdout_logfile=/var/log/supervisor/buykori-worker.out.log
stderr_logfile=/var/log/supervisor/buykori-worker.err.log
stdout_logfile_maxbytes=10MB
stdout_logfile_backups=5
stderr_logfile_maxbytes=10MB
stderr_logfile_backups=5
environment=PATH="/var/www/buykori-adsync/venv/bin:%(ENV_PATH)s",EVENT_INGEST_MODE="redis_stream",USAGE_RESERVATION_MODE="worker",DB_POOL_SIZE="5",DB_MAX_OVERFLOW="3",HTTP_TIMEOUT_SECONDS="8",EVENT_WORKER_BATCH_SIZE="50",EVENT_REDIS_STREAM_BATCH_SIZE="50"
